# arc/prepm2ubuntu.conf
start on filesystem and started networking
respawn
chdir /root/deploy/yebo
exec /usr/local/bin/node server.js